# deploy/poly-executor.service
[Unit]
Description=Poly Executor Dublin (quotes maker Polymarket)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rust-quant-bot-v2/backend
EnvironmentFile=/home/ubuntu/rust-quant-bot-v2/backend/.env
ExecStart=/home/ubuntu/rust-quant-bot-v2/backend/target/release/polymarket_mm_bot
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
