# deploy/rust-quant-bot-radar.service
[Unit]
Description=Rust Quant Bot — RADAR (Tokyo, ap-northeast-1) : OBI Binance+OKX → UDP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rust-quant-bot
# IP publique de l'exécuteur (Dublin) — À RENSEIGNER au déploiement.
Environment=TARGET_EXECUTOR_IP=REMPLACER_PAR_IP_PUBLIQUE_DUBLIN
Environment=TARGET_PORT=8080
# Dashboard local du radar (OBI/vélocité/signaux émis + latences CEX).
Environment=PORT=8768
Environment=RUST_LOG=info
ExecStart=/home/ubuntu/rust-quant-bot/target/release/rust-quant-bot radar
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/rust-quant-bot/data/radar.log
StandardError=append:/home/ubuntu/rust-quant-bot/data/radar.log

[Install]
WantedBy=multi-user.target
