# deploy/rust-quant-bot-radar.service
[Unit]
Description=Rust Quant Bot — RADAR (Tokyo, ap-northeast-1) : OBI Binance+OKX → UDP (live + paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rust-quant-bot
# Cible LIVE (Dublin) — servie EN PREMIER. À RENSEIGNER au déploiement.
Environment=TARGET_LIVE_IP=REMPLACER_PAR_IP_PUBLIQUE_DUBLIN
Environment=TARGET_LIVE_PORT=8080
# Cible PAPER (machine séparée) — tir SECONDAIRE (après le live). Optionnel : si absent, live seul.
Environment=TARGET_PAPER_IP=REMPLACER_PAR_IP_PUBLIQUE_PAPER
Environment=TARGET_PAPER_PORT=8081
# Dashboard local du radar (OBI/vélocité/signaux émis + latences CEX).
Environment=PORT=8768
Environment=RUST_LOG=info
ExecStart=/home/ubuntu/rust-quant-bot/target/release/rust-quant-bot radar
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/rust-quant-bot/data/radar.log
StandardError=append:/home/ubuntu/rust-quant-bot/data/radar.log

[Install]
WantedBy=multi-user.target
