# deploy/rust-quant-bot-executor.service
[Unit]
Description=Rust Quant Bot — EXÉCUTEUR (Dublin, eu-west-1) : écoute UDP → Polymarket (paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rust-quant-bot
# PAPER obligatoire tant que le signing EIP-712 n'existe pas (règle de sécurité du projet).
Environment=DRY_RUN=true
Environment=LISTEN_PORT=8080
# Dashboard local de l'exécuteur (PM/position/PnL).
Environment=PORT=8768
Environment=STATE_PATH=/home/ubuntu/rust-quant-bot/data/sniper_state.json
Environment=TRADES_PATH=/home/ubuntu/rust-quant-bot/data/sniper_trades.jsonl
Environment=START_CASH=200
Environment=RUST_LOG=info
ExecStart=/home/ubuntu/rust-quant-bot/target/release/rust-quant-bot executor
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/rust-quant-bot/data/executor.log
StandardError=append:/home/ubuntu/rust-quant-bot/data/executor.log

[Install]
WantedBy=multi-user.target
